# templates/nodejs/nodejs-supervisord.conf
[supervisord]
nodaemon=true

[program:permissions]
command=chmod +x /root/templates/nodejs-init.conf

[program:init]
command=bash /root/templates/nodejs-init.conf


