# app/initializer
#!/sbin/openrc-run

start() {
    ebegin "Starting GorillaJS Guest"
    start-stop-daemon --start --background \
    --exec forever /home/gorilla/gorillajs/app/index.js guest \
    -u gorilla
    eend $?
}
