# templates/wordpress/mysql-supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command="./root/templates/mysql-init.conf"
