# templates/common/mysql-supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=bash /root/templates/mysql-init.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

