# templates/common/mysql-supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=bash /root/templates/mysql-init.conf

