# templates/blank+database/apache-supervisord.conf
[supervisord]
nodaemon=true

[program:permissions]
command=chmod +x /root/templates/apache-init.conf

[program:init]
command=bash root/templates/apache-init.conf


